# services/bloom-botanics.service
[Unit]
Description=BloomBotanics Agricultural Monitoring System
Documentation=https://github.com/yourusername/BloomBotanics
After=network-online.target multi-user.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/BloomBotanics
Environment=PATH=/home/pi/BloomBotanics/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/pi/BloomBotanics
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/BloomBotanics/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bloom-botanics

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/pi/BloomBotanics/data

[Install]
WantedBy=multi-user.target
